# apps/workspace/deploy/ai.hermes.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>ai.hermes.dashboard</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/phill-mac/.hermes/hermes-agent/venv/bin/python</string>
		<string>-m</string>
		<string>hermes_cli.main</string>
		<string>dashboard</string>
		<string>--no-open</string>
		<string>--port</string>
		<string>9119</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/phill-mac/.hermes/logs/dashboard.error.log</string>
	<key>StandardOutPath</key>
	<string>/Users/phill-mac/.hermes/logs/dashboard.log</string>
</dict>
</plist>
