# apps/workspace/deploy/ai.hermes.workspace.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>ai.hermes.workspace</string>
  <key>ProgramArguments</key><array>
    <string>/bin/bash</string>
    <string>/Users/phill-mac/pi-seo-workspace/unite-group/apps/workspace/start-operator.sh</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/phill-mac/.hermes/logs/workspace.log</string>
  <key>StandardErrorPath</key><string>/Users/phill-mac/.hermes/logs/workspace.error.log</string>
</dict></plist>
